# hey-marketing-reply/reply.service
[Unit]
After=network.target
[Service]
User=<your-user>
ExecStart=/path/to/hey-marketing-reply/reply.sh
Restart=on-failure
[Install]
WantedBy=multi-user.target
